# systemd_units/pinenote_dbus_service.service
[Unit]
Description=Pinenote DBUS service

[Service]
ExecStart=/usr/bin/pinenote_dbus_service

[Install]
WantedBy=multi-user.target
